# Documents/systemd/backup.service
[Unit]
Description=Backup to raspberry pi

[Service]
Type=simple
User=root
ExecStart=/home/philipj/Documents/systemd/backup.sh

